# roles/nginx_conf_rsync/files/nginx_syncfile.conf
[program:nginx-sync]
command=/usr/local/pyenv/versions/env352/bin/python syncfile.py
directory=/root/nginx_conf_rsync
user=root
